# ubuntu_utopic_registry/registry_supervisord.conf
[program:registry]
command=docker-registry
